# vm-assets/systemd/docker.service.d/storage.conf
[Unit]
Requires=init-var.service
After=init-var.service

[Service]
StandardOutput=journal+console
StandardError=journal+console
